# pkg/DEBIAN/postinst
#!/bin/bash

PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

case $1 in
  abort-upgrade)
    ;;
  configure|abort-remove)
    systemctl preset flecs
    systemctl is-enabled flecs >/dev/null && systemctl start flecs
    systemctl is-enabled flecs-webapp >/dev/null 2>&1 && systemctl start flecs-webapp
    ;;
esac

exit 0
