# backend/pkg/debian/DEBIAN/postinst
#!/bin/bash

PATH=/sbin:/bin:/usr/sbin:/usr/bin

case "$1" in
  *)
    docker network create flecs-internal >/dev/null 2>&1
    systemctl enable --now ##PACKAGE##
    ;;
esac

exit 0
