# updater/pkg/debian/DEBIAN/prerm
#!/bin/bash

PATH=/sbin:/bin:/usr/sbin:/usr/bin

case "$1" in
  remove)
      systemctl disable ##PACKAGE##
    ;;
  *)
    ;;
esac

exit 0
